# scripts/ccr.service
[Unit]
Description=Claude Code Router (musistudio/claude-code-router)
Documentation=https://github.com/musistudio/claude-code-router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PATH=/home/admin/.nvm/versions/node/v24.15.0/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/admin/.nvm/versions/node/v24.15.0/bin/ccr start
ExecStop=/home/admin/.nvm/versions/node/v24.15.0/bin/ccr stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
